# units/darktheme.service
[Unit]
Description=Enable KDE breeze dark theme

[Service]
Type=simple
ExecStartPre=/bin/sleep 2
ExecStart=/usr/local/bin/theme dark
