# units/darktheme.service
[Unit]
Description=Enable KDE breeze dark theme

[Service]
Type=simple
ExecStart=/usr/local/bin/darktheme
