# units/lighttheme.service
[Unit]
Description=Enable KDE breeze (light) theme

[Service]
Type=simple
ExecStartPre=/bin/sleep 5
ExecStart=/usr/local/bin/theme light
